# qwen-local-delegate-mcp.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>local.qwen-local-delegate-mcp</string>

    <key>ProgramArguments</key>
    <array>
        <string>__UV_PATH__</string>
        <string>run</string>
        <string>--script</string>
        <string>__REPO_DIR__/src/server.py</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>MCP_TRANSPORT</key>
        <string>streamable-http</string>
        <key>MCP_HOST</key>
        <string>127.0.0.1</string>
        <key>MCP_PORT</key>
        <string>11435</string>
        <key>QWEN_MODEL</key>
        <string>__QWEN_MODEL__</string>
        <key>PATH</key>
        <string>__PATH__</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>__REPO_DIR__</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>__REPO_DIR__/data/server.log</string>

    <key>StandardErrorPath</key>
    <string>__REPO_DIR__/data/server.log</string>
</dict>
</plist>
